# ble-scanner.service
[Unit]
Description=BLE Scanner
After=network-online.target bluetooth.target
Wants=network-online.target

[Service]
User=qatester
WorkingDirectory=/home/qatester/ble_scanner
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStartPre=/bin/sh -c 'for i in 1 2 3 4 5; do hciconfig hci0 up 2>/dev/null && break; sleep 2; done'
ExecStart=/usr/bin/python3 main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
